# NBitcoinTest/crontab
*/15 * * * * root /usr/bin/dotnet app/NBitcoinTest.dll